# scripts/templates/app.service.in
[Unit]
Description=Thermo Logger Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
Group=root

WorkingDirectory=@APP_INSTALL_ROOT@
ExecStart=node dist/main.js --port=80

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
